# examples/systemd/spool-tag-writer-agent.service
[Unit]
Description=Spool Tag Writer Moonraker Agent
After=network-online.target moonraker.service
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/spool-tag-writer
EnvironmentFile=/home/pi/spool-tag-writer/.env
ExecStart=/home/pi/spool-tag-writer/.venv/bin/python -m app.moonraker_agent
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
